# hyprsuite/.config/systemd/user/hyprlock-on-resume.service
[Unit]
Description=Run hyprlock on resume
After=supend.target

[Service]
ExecStart=/usr/bin/hyprlock
Type=oneshot

[Install]
WantedBy=suspend.target
